# roles/synchronizer/scripts/launchd/com.exocortex.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exocortex.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>__IWE_WORKSPACE__/DS-strategy/roles/synchronizer/scripts/scheduler.sh</string>
        <string>dispatch</string>
    </array>

    <!-- Пробуждение в ключевые часы:
         0:00  — week-review (Пн) + code-scan
         3:00  — extractor inbox-check
         4:00  — strategist morning
         6:00  — extractor + daily-report
         9:00  — extractor
         12:00 — extractor
         15:00 — extractor
         18:00 — extractor
         21:00 — extractor
         23:00 — strategist note-review
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/logs/synchronizer/launchd-scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/logs/synchronizer/launchd-scheduler-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>__USER_HOME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
